# monero-miner.service
[Unit]
Description=Monero XMR Mining Service
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/mining/folder
ExecStart=/usr/bin/python3 /path/to/mining/folder/main.py
Restart=always
RestartSec=10
StandardOutput=append:/path/to/mining/folder/mining.log
StandardError=append:/path/to/mining/folder/mining.log

# Optional: Resource limits
# LimitNOFILE=65536
# Nice=10

[Install]
WantedBy=multi-user.target
